# systemd/easyape.service
[Unit]
Description=EasyApe - Text to Stake
After=network.target

[Service]
Type=simple
User=root

# ✅ IMPORTANT: Installer should replace this with actual repo path
WorkingDirectory=__EASYAPE_ROOT__

ExecStart=__EASYAPE_ROOT__/.venv/bin/python \
    -u -m stakechat_bot.main \
    --config __EASYAPE_ROOT__/config.yaml

Restart=always
RestartSec=5

# Cleaner shutdown behavior
TimeoutStopSec=10
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
